# etopai/.deb/prerm
#!/bin/sh

case "$1" in
	purge)
		chown -R root:root /var/lib/etopa
        deluser etopa
		;;
esac

systemctl stop etopa.service
systemctl disable etopa.service

exit 0
