# etopai/.deb/prerm
#!/bin/sh

chown -R root:root /var/lib/etopa
deluser etopa
systemctl stop etopa.service
systemctl disable etopa.service

exit 0
