# etopai/.deb/prerm
#!/bin/sh

systemctl stop etopa.service

if [ $1 = "remove" ] || [ $1 = "purge" ]; then
  systemctl disable etopa.service
  chown -R root:root /var/lib/etopa
  deluser etopa
fi

exit 0
